# ansible/roles/restic_backup/templates/restic-backup.timer.j2
# {{ ansible_managed }}
[Unit]
Description=Run Restic backup daily

[Timer]
OnCalendar=*-*-* {{ restic_timer }}:00
Persistent=true

[Install]
WantedBy=timers.target